# root/etc/services.d/samba/run
#!/usr/bin/with-contenv bash

/usr/sbin/smbd --foreground --log-stdout

